# .config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=ruby /home/antoine/.conky/starter.rb
Comment=Custom conky desktop
